# packer/webapp.path
[Unit]
Description=Watch for the appearance of environment file

[Path]
PathExists=/etc/environment
Unit=webapp.service

[Install]
WantedBy=multi-user.target