# packer/webapp.path
[Unit]
Description=Watch for the appearance of .env file

[Path]
PathExists=/opt/.env
Unit=webapp.service

[Install]
WantedBy=multi-user.target